# django/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command =uwsgi --ini  /var/www/app/uwsgi.ini
